# deploy/local-ai-llm.service.example
[Unit]
Description=Local AI LLM Monitor and Portal
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/local-ai-llm
EnvironmentFile=-/opt/local-ai-llm/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5
User=local-ai-llm
Group=local-ai-llm
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/local-ai-llm/config

[Install]
WantedBy=multi-user.target
